# debian/pleco.cron.d
* * * * *   root /etc/NetworkManager/dispatcher.d/10-push-ip
